# install/launchagents/com.nelson.comms-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.nelson.comms-agent</string>

  <!-- Starts the comms agent (detached tmux running `claude --remote-control comms`)
       at login. launch.sh is a no-op if the session is already running.
       No KeepAlive: launch.sh exits immediately after spawning tmux, and
       resurrection-on-death is the Tickle supervisor's job (phase 2), not launchd's. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nelson/agents/comms/launch.sh</string>
  </array>

  <key>RunAtLoad</key><true/>

  <key>StandardOutPath</key><string>/Users/nelson/.local/state/comms-agent/launchd.log</string>
  <key>StandardErrorPath</key><string>/Users/nelson/.local/state/comms-agent/launchd.log</string>
</dict>
</plist>
