# install/launchagents/com.nelson.refresh-contacts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nelson.refresh-contacts</string>

    <!-- Rebuilds ~/obsidian/.obsidian/contacts.json from macOS AddressBook -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nelson/.local/bin/refresh-contacts</string>
        <string>--quiet</string>
    </array>

    <!-- Daily at 03:55 — 5 min before com.nelson.legal-mail.daily (04:00),
         so canonicalize.py reads a fresh contacts.json on the pull run. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>

    <!-- Don't fire on launchctl load — wait for the scheduled time -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logs: tail with `tail -F ~/Library/Logs/refresh-contacts.log` -->
    <key>StandardOutPath</key>
    <string>/Users/nelson/Library/Logs/refresh-contacts.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nelson/Library/Logs/refresh-contacts.log</string>

    <!-- Standard PATH so /usr/bin/python3 etc. resolve -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- launchd will skip the run if the previous one is still going -->
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
