# install/launchagents/dev.tickle.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.tickle.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nelson/Library/Application Support/tickle/bin/tickle</string>
    <string>daemon</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Jobs/scripts are versioned in dotfiles (mirrored across machines);
         runtime state (runs/state/logs/bin) stays per-machine in the default
         DATA_HOME (~/Library/Application Support/tickle). -->
    <key>TICKLE_CONFIG_HOME</key>
    <string>/Users/nelson/repos/dotfiles/tickle</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/nelson/Library/Application Support/tickle/logs/daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nelson/Library/Application Support/tickle/logs/daemon.err.log</string>
</dict>
</plist>
